# docs/web-scheduler.service.example
[Unit]
Description=Web Scheduler API (Uvicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/web_scheduler
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/web_scheduler/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target


